# templates/shutup-gnome.service.j2
[Unit]
Description=Ebay Shipping Last Call Notifier

[Service]
Type=oneshot
ExecStart={{ node }} {{ entrypoint }}
Environment="DISPLAY=:0"
User={{ user }}
Group={{ group }}
