# templates/shutup-gnome.timer.j2
[Unit]
Description=Shutup gnome

[Timer]
OnCalendar={{ calendar }}

[Install]
WantedBy=multi-user.target
